# vagrant/supervisor/envelope.conf
[program:envelope]
command = /vagrant/bin/start_gunicorn.sh
user = root
stdout_logfile = /var/log/supervisor/envelope.log
redirect_stderr = true